# infra/systemd/ghs-db-backup-daily.service
[Unit]
Description=GHS PostgreSQL Daily Backup
After=network.target postgresql.service
Wants=ghs-db-backup-daily.timer

[Service]
Type=oneshot
User=root
Environment="DB_NAME=ghs_db"
Environment="DB_USER=ghs_db_super_user"
Environment="DB_HOST=127.0.0.1"
Environment="DB_PORT=5432"
Environment="BACKUP_DIR=/var/backups/ghs/postgres"
Environment="BACKUP_RETENTION_DAILY=30"
ExecStart=/usr/local/bin/ghs-pg-backup --daily
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ghs-db-backup-daily

[Install]
WantedBy=multi-user.target
