# infra/systemd/ghs-web-green.service
[Unit]
Description=GHS Web Server (Green Slot)
After=network.target
Wants=ghs-web-green.service

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/apps/ghs-green/apps/web
Environment="NODE_ENV=production"
Environment="SLOT=green"
Environment="PORT=5176"
EnvironmentFile=/etc/ghs/.env.production
ExecStart=/usr/bin/node --loader vite/esm src/server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ghs-web-green

[Install]
WantedBy=multi-user.target
