# infra/systemd/ghs-db-backup.service
[Unit]
Description=GHS PostgreSQL daily backup
After=postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=postgres
Environment=DB_NAME=ghs_db
Environment=DB_USER=ghs_db_super_user
Environment=BACKUP_DIR=/var/backups/ghs/postgres
Environment=RETAIN_DAYS=7
ExecStart=/usr/local/bin/ghs-pg-backup
